# autostart/custom.shell-on-startup.desktop
[Desktop Entry]
Terminal=false
Name=startup shell
Type=Application
Exec=alacritty
Icon=shell
StartupNotify=false
