# autostart/custom.polybar.desktop
[Desktop Entry]
Terminal=false
Name=polybar
Type=Application
Exec=polybar
Icon=taskbar
NoDisplay=false
StartupNotify=false
